# risc0-asteroids-verifier/deploy/systemd/risc0-asteroids-api.service
# NOTE: Update WorkingDirectory and ExecStart below to match your actual
# clone path.  The VASTAI script defaults to /workspace/stellar-zk, but a
# manual `git clone` of the repo will use /workspace/stellar-zk-claude.

[Unit]
Description=RISC0 Asteroids Verifier API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/workspace/stellar-zk/risc0-asteroids-verifier
EnvironmentFile=/etc/stellar-zk/api-server.env
ExecStart=/workspace/stellar-zk/risc0-asteroids-verifier/target/release/api-server
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
